# setup/check4sound/soundcheck.conf
[program:check4sound]
command=/home/dietpi/.venv/bin/python /home/dietpi/check4sound.py
autostart=true
autorestart=true
priority=10
