# setup/check4sound/soundcheck.conf
[program:check4sound]
command=/home/tom/.venv/bin/python /home/tom/check4sound.py
autostart=true
autorestart=true
priority=10
